# services/subsonic/run
#!/bin/bash
echo "Starting subsonic"
/etc/init.d/subsonic start
tail -f /config/subsonic.log
